# deploy/finance-backup.timer
# Runs rendered finance-backup.service daily (and catches up if the machine was
# off). Install this timer only with the two rendered and verified services.
[Unit]
Description=Daily Finance Engine backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
